# backend/crontab
# NBA Stat Spot Scheduled Tasks
# All times are in UTC

# Refresh daily props cache every day at 6:00 AM UTC (1-2 AM EST)
# This ensures fresh data is available for the day's games
0 6 * * * curl -f -X POST "http://localhost:8000/api/v1/admin/refresh/daily-props?min_confidence=50&limit=50" || echo "Failed to refresh daily props"

# Refresh high hit rate cache every day at 6:05 AM UTC
5 6 * * * curl -f -X POST "http://localhost:8000/api/v1/admin/refresh/high-hit-rate?min_hit_rate=0.75&limit=10&last_n=10" || echo "Failed to refresh high hit rate"

# Cleanup expired cache entries daily at 3:00 AM UTC (10-11 PM EST previous day)
# This removes old cache entries to keep the database clean
0 3 * * * curl -f -X POST "http://localhost:8000/api/v1/admin/cache/cleanup" || echo "Failed to cleanup expired cache"

# Optional: Refresh all caches weekly on Sunday at 2:00 AM UTC
# This ensures a full refresh once per week
# 0 2 * * 0 curl -f -X POST "http://localhost:8000/api/v1/admin/refresh/all" || echo "Failed to refresh all caches"

